# mediaelch/rootfs/etc/services.d/openbox/run
#!/command/with-contenv bashio

export DISPLAY=":0"

bashio::log.info "Starting OpenBox...."
exec /usr/bin/openbox-session
