# rootfs/etc/services.d/fail2ban/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: fail2ban
# Runs the fail2ban
# ==============================================================================

bashio::log.info "Starting fail2ban..."

# Run
exec fail2ban-server -f -x -v start
