# hpssacli/rootfs/etc/s6-overlay/s6-rc.d/mqtt-handler/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start mqtt service
# ==============================================================================

chmod -R 755 /media/runtest
/media/runtest
	