# fail2ban/rootfs/etc/services.d/fail2ban/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: fail2ban
# Take down the S6 supervision tree when fail2ban fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "fail2ban crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "fail2ban stopped, restarting..."